# ansible/roles/dev-desktop/templates/get-all-crates.timer.j2
[Unit]
Description=Refresh crates.io export timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=3600

[Install]
#WantedBy=timers.target
