# ansible/roles/letsencrypt/templates/renew-ssl-certs.timer
[Unit]
Description=Renew SSL certificates each week

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
