# ansible/roles/playground/templates/playground.slice
[Unit]
Description=Resource management group for playground processes
Before=slices.target
